# packages/configs-core/profile.d/10-tmpdir.sh
#!/bin/sh

export TMPDIR="$HOME/tmp"
mkdir -p "$TMPDIR"
